# dynwg.service
[Unit]
Description=WireGuard DynDNS watchdog
Requires=systemd-networkd.service

[Service]
Type=oneshot
ExecStart=/usr/bin/dynwg